# src/gps-share/debian/gps-share.service
[Unit]
Description=Utility to share your GPS device on local network
User=gps-share
Requires=gps-share.socket

[Service]
ExecStart=/usr/bin/gps-share --socket-path /run/gps-share.sock
Restart=on-failure

[Install]
WantedBy=default.target
